# systemd/thehive-consumer.service
[Unit]
Description=SOCIA Kafka to TheHive alert consumer
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=socia-thehive
Group=socia-thehive
WorkingDirectory=/opt/socia-thehive/consumer
EnvironmentFile=/opt/socia-thehive/consumer/.env
ExecStart=/opt/socia-thehive/consumer/venv/bin/python /opt/socia-thehive/consumer/thehive-consumer.py
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
